# data/ayatana-indicator-printers.service.in
[Unit]
Description=Ayatana Indicator Printers Service
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/ayatana-indicator-printers/ayatana-indicator-printers-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
