# data/indicator-printers.service.in
[Unit]
Description=Indicator Printers Service
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/indicator-printers-service
